# Praktikum 2025/D03/praktikum-modul-3-d03-2025/osboot/myramdisk/etc/inittab
TO:23:respawn:/bin/getty -L -ttyS0 115200 vt100
