# Praktikum 2025/D03/praktikum-modul-1-d03-2025/task-3/crontabs
*dijakankan di awal*
*/2 * * * * /home/freddskiii/automation.sh

*dalam automation.sh*
* */2 * * * /home/freddskiii/archive.sh
*/10 * * * * /home/freddskiii/download.sh
